# templates/etc/systemd/system/certbot-renew-certbot_cert_name.service.j2
[Unit]
Description=Renew Let's Encrypt certificates for {{ certbot_cert_name }}

[Service]
Type=oneshot
# check for renewal, only start/stop nginx if certs need to be renewed
ExecStart={{ certbot_autorenew_command }}
